# riemann/debian/riemann.upstart
description "riemann server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 60
normal exit 0 2

limit nofile 8192 8192
setuid _riemann
setgid _riemann

pre-start script
    [ -r "/usr/share/java/riemann.jar" ] || exit 0
    [ -r "/etc/riemann/riemann.config" ] || exit 0
end script

script
    . /etc/default/riemann
    exec $JAVA $JAVA_OPTS -cp /usr/share/java/riemann.jar:$EXTRA_CLASSPATH riemann.bin start /etc/riemann/riemann.config
end script
