# puppetboard/debian/puppetboard.upstart
description "puppetboard"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 60
normal exit 0 2

limit nofile 8192 8192
setuid _puppetboard
setgid _puppetboard
console log

env LISTEN=127.0.0.1:9090
env PUPPETBOARD_SETTINGS=/etc/puppetboard/config.py

pre-start script
    [ -r $PUPPETBOARD_SETTINGS ] || exit 0
end script

script
    exec /usr/share/python/puppetboard/bin/gunicorn -b $LISTEN puppetboard.app:app
end script
